# grott-complete/rootfs/etc/s6-overlay/s6-rc.d/grott-server/run
#!/usr/bin/with-contenv bashio

# Pass verbose setting to grottserver.py via env var
VERBOSE=$(bashio::config 'verbose')
export GROTTSERVER_VERBOSE="${VERBOSE}"

bashio::log.info "Starting Grott server API on port 5782 (verbose=${VERBOSE})..."
cd /app

# When not verbose, redirect stdout to /dev/null (keeps stderr for errors)
if [ "${VERBOSE}" = "true" ]; then
    exec python3 -u grottserver.py
else
    exec python3 -u grottserver.py > /dev/null
fi
