# grott-complete/rootfs/etc/s6-overlay/s6-rc.d/grott-server/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Grott server API on port 5782..."
cd /app
exec python3 -u grottserver.py
