# grott-complete/rootfs/etc/s6-overlay/s6-rc.d/grott-proxy/run
#!/usr/bin/with-contenv bashio

VERBOSE=$(bashio::config 'verbose')

bashio::log.info "Starting Grott proxy on port 5279 (verbose=${VERBOSE})..."
cd /app

# When not verbose, redirect stdout to /dev/null (keeps stderr for errors)
if [ "${VERBOSE}" = "true" ]; then
    exec python3 -u grott.py -c /app/grott.ini
else
    exec python3 -u grott.py -c /app/grott.ini > /dev/null
fi
