# grott-complete/rootfs/etc/s6-overlay/s6-rc.d/grott-proxy/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Grott proxy on port 5279..."
cd /app
exec python3 -u grott.py -c /app/grott.ini
